# configs/systemd-services/eupnea-update.service
[Unit]
Description=Eupnea auto-update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-scripts